# snmpd-extend/snmpd-extend/debian/postrm
#!/bin/sh
# postrm script for snmpd-extend
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postrm> `remove'
#        * <postrm> `purge'
#        * <old-postrm> `upgrade' <new-version>
#        * <new-postrm> `failed-upgrade' <old-version>
#        * <new-postrm> `abort-install'
#        * <new-postrm> `abort-install' <old-version>
#        * <new-postrm> `abort-upgrade' <old-version>
#        * <disappearer's-postrm> `disappear' <overwriter>
#          <overwriter-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

export PATH="/usr/sbin:/usr/bin:/sbin:/bin"
export LC_ALL="C"

SNMP_CONF="/etc/snmp/snmpd.local.conf"
BEGIN_DELIMITER="### Begin snmpd-extend ###"
END_DELIMITER="### End snmpd-extend ###"
CODENAME=$(lsb_release --codename --short)

is_old_distribution () {
    if [ "$CODENAME" = "etch" ] || [ "$CODENAME" = "lenny" ]
    then
        true
    else
        false
    fi
}

if is_old_distribution
then
    SUDOERS_CONF="/etc/sudoers"
else
    SUDOERS_CONF="/etc/sudoers.d/snmpd-extend"
fi




case "$1" in
    purge|remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
        if is_old_distribution
        then
            tmp_file=$(mktemp)

            # We remove the "snmpd-extend" bloc (if exists) in the "sudoers" file.
            sed -n '/^'"$BEGIN_DELIMITER"'$/,/^'"$END_DELIMITER"'$/!p' "$SUDOERS_CONF" > "$tmp_file"

            # We remove the empty lines (if exist) in $tmp_file.
            sed -ni '1h;1!H;${x;s/[[:space:]]*$//;p}' "$tmp_file"

            # We add just one empty line at the end.
            echo "" >> "$tmp_file"

            cat "$tmp_file" > "$SUDOERS_CONF"
            rm -f "$tmp_file"
        else
            if [ -e "$SUDOERS_CONF" ]; then
                rm -f "$SUDOERS_CONF"
            fi
        fi

        # We remove the snmpd.local.conf file (if exist).
        rm -f "$SNMP_CONF"

        if [ "$1" = "remove" ]
        then
            # We just want to 1) stop the daemon snmpd before
            # install/remove/purge/upgrade and 2) start after.
            # See the diagrams here https://wiki.debian.org/MaintainerScripts
            if [ -x "/etc/init.d/snmpd" ]; then
                echo "Starting snmpd..."
                sleep 1
                invoke-rc.d snmpd start || true
            fi
        fi
    ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
