# shinken-packs/shinken-packs/debian/postinst
#!/bin/sh
# postinst script for shinken-packs
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

LIBEXEC_DIR="/usr/share/shinken-packs/libexec"
CHECK_ICMP="/usr/lib/nagios/plugins/check_icmp"
VERSION='__VERSION__'

case "$1" in
    configure)
        for f in "$LIBEXEC_DIR/check_"*; do
            [ "$f" = "$LIBEXEC_DIR/check_*" ] && continue
            chown root:root "$f"
            chmod 755 "$f"
        done

        if [ -e "$CHECK_ICMP" ]; then
            chmod u+s "$CHECK_ICMP" || true
        fi

        # Starting of shinken only after a "real" upgrade.
        if [ -x "/etc/init.d/shinken" ] && [ -n "$2" ]; then
            if [ "$2" != "$VERSION" ]; then
                # This is a "real" upgrade, not a new install (after a remove).
                echo "Starting shinken after an upgrade..."
                sleep 1
                invoke-rc.d shinken start || true
            fi
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
